# build/templates/nixos/nixos-vm-template/nixos-flake-update.service
[Unit]
Description=Update NixOS flake and rebuild system

[Service]
Type=oneshot
ExecStart=/etc/nixos/nixos-flake-update.sh 